# dist/gae/deploy/saycbridge-upstart.conf
description "service user saycbridge upstart daemon"

start on runlevel [2345]
respawn
setuid saycbridge

script
  export HOME=/home/saycbridge
  export XDG_RUNTIME_DIR="$HOME/.cache/run"
  mkdir -p $XDG_RUNTIME_DIR
  chmod 0700 $XDG_RUNTIME_DIR

  exec /sbin/init --user
end script
