# dist/gae/deploy/saycbridge.conf
description "saycbridge app server"

start on started
respawn
kill signal SIGINT

script
  cd /home/saycbridge
  . virtualenv/bin/activate
  cd deploy/dist/gae
  PYTHONPATH=/home/saycbridge/z3/build exec python production_main.py
end script
