# dist/gae/deploy/saycbridge.conf
description "saycbridge app server"

start on runlevel [2345]
#manual
#expect default
setuid saycbridge
respawn
kill signal SIGINT

script
  cd /home/saycbridge/deploy
  . sayc-env/bin/activate
  cd dist/gae
  PYTHONPATH=/home/saycbridge/z3/build python production_main.py
end script
