# dist/gae/deploy/saycbridge.conf
description "saycbridge app server"

start on started
respawn
kill signal SIGINT

script
  cd /home/saycbridge
  . virtualenv/bin/activate
  cd deploy
  make serve-prod
end script
