# dist/gae/deploy/saycbridge.conf
description "saycbridge app server"

start on started
respawn
kill signal SIGINT

script
  ulimit -m 400000 # Because saycbridge can't be trusted not to use all memory!
  cd /home/saycbridge
  . virtualenv/bin/activate
  cd deploy
  make serve-prod
end script
